# chezmoi/dot_config/autostart/discord-stable.desktop
[Desktop Entry]
Type=Application
Exec=discord
Name=Discord
Icon=discord
Comment=Text and voice chat for gamers.
X-GNOME-Autostart-enabled=true
