# chezmoi/dot_config/autostart/private_discord.desktop
[Desktop Entry]
Categories=Network;InstantMessaging
Comment[en_US]=
Comment=
Exec=Discord --start-minimized
GenericName[en_US]=All-in-one cross-platform voice and text chat for gamers
GenericName=All-in-one cross-platform voice and text chat for gamers
Icon=discord
MimeType=
Name[en_US]=Discord
Name=Discord
Path=
StartupNotify=true
StartupWMClass=discord
Terminal=false
TerminalOptions=
Type=Application
Version=1.4
X-KDE-SubstituteUID=false
X-KDE-Username=
