# chezmoi/dot_config/autostart/udiskie.desktop
[Desktop Entry]
Type=Application
Name=udiskie
Comment=udiskie
Exec=udiskie
StartupNotify=false
Terminal=false
